# deploy/systemd/mb-deploy.service
[Unit]
Description=Material Bank — auto-deploy from GitHub main (CD)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/env bash __APP_DIR__/deploy/auto_deploy.sh
TimeoutStartSec=900
StandardOutput=append:__APP_DIR__/reports/mb-deploy.log
StandardError=append:__APP_DIR__/reports/mb-deploy.log
