# deploy/systemd/mb-harvest.service
[Unit]
Description=Material Bank — one harvest sweep (triggered hourly by mb-harvest.timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=-__APP_DIR__/.env
ExecStart=/usr/bin/env bash __APP_DIR__/deploy/harvest_once.sh
# a sweep draining a giant can run long; don't let systemd kill it
TimeoutStartSec=infinity
StandardOutput=append:__APP_DIR__/reports/mb-harvest.log
StandardError=append:__APP_DIR__/reports/mb-harvest.log
