# deploy/systemd/mb-harvest.service
[Unit]
Description=Material Bank — one harvest sweep (triggered hourly by mb-harvest.timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=-__APP_DIR__/.env
ExecStart=/usr/bin/env bash __APP_DIR__/deploy/harvest_once.sh
# The queue-drain phase is time-boxed inside the script (resumable next sweep).
# This is only a backstop so a wedged sweep can never monopolize the timer the
# way an unbounded one did (a 2-day run blocked all recovery/bespoke/self-heal).
TimeoutStartSec=2h
TimeoutStopSec=90s
KillMode=mixed
StandardOutput=append:__APP_DIR__/reports/mb-harvest.log
StandardError=append:__APP_DIR__/reports/mb-harvest.log
